# appliance/mkosi.extra/etc/systemd/system/getty@tty1.service.d/spatium-console-gate.conf
# #393 — gate getty@tty1 on the appliance console mode.
#
# The Talos console dashboard (spatium-console.service) is a getty
# REPLACEMENT on tty1: it runs when the kernel cmdline does NOT carry
# spatium-console=off (its ConditionKernelCommandLine=!spatium-console=off),
# and Conflicts=getty@tty1.service. In the "text_console" console mode the
# grub menuentry passes spatium-console=off, so the dashboard's condition
# fails and it doesn't start — but pre-#393 NOTHING claimed tty1 then, so an
# operator who picked the plain-console mode was left at bare kernel text with
# NO login prompt at all.
#
# This drop-in makes getty@tty1 the mirror image of the dashboard: it runs
# ONLY when spatium-console=off is present. Combined with the
# getty.target.wants symlink (added in mkosi.postinst), that means:
#   * dashboard / verbose_dashboard mode (no spatium-console=off) → this
#     condition fails → getty@tty1 stays inactive → the dashboard owns tty1.
#   * text_console mode (spatium-console=off) → the dashboard's condition
#     fails, THIS condition passes → getty@tty1 starts → a normal login.
# The two are mutually exclusive by condition, so they never fight over tty1.
#
# tty2..tty6 are unaffected (this is the tty1 instance drop-in only), so the
# Alt+F2..F6 shell escape hatch via logind autovt still works in every mode.
[Unit]
ConditionKernelCommandLine=spatium-console=off
