# appliance/mkosi.extra/etc/systemd/system/getty@tty1.service.d/spatium-mode-skip.conf
[Unit]
# Skip starting agetty on tty1 when the kernel cmdline carries
# spatium-mode=install or spatium-mode=live. In those modes the
# spatium-install.service / spatium-live-banner.service unit owns
# tty1 and runs whiptail there.
#
# Why a condition drop-in instead of Conflicts= on the spatium
# units: Conflicts= is registered at unit-load time, so even
# though the spatium units have ConditionKernelCommandLine= that
# prevents activation on the *installed* system (no spatium-mode=),
# the conflict relationship is still recorded — getty@tty1 ends up
# in ConflictedBy with units it will never actually run with, and
# never auto-starts at multi-user.target on the installed system.
#
# Reversing the gate (condition on getty@tty1 itself) keeps the
# same end result without the unit-load-time leak: getty@tty1 only
# attempts to start when neither spatium-mode is in cmdline. On
# the live ISO it's a no-op (condition fails), no race with the
# wizard, no autologin-flash. On the installed system the
# condition passes and getty@tty1 runs as Debian intended.
ConditionKernelCommandLine=!spatium-mode=install
ConditionKernelCommandLine=!spatium-mode=live
