# appliance/mkosi.extra/etc/systemd/system/root.mount
[Unit]
Description=/root (bind from /var/root — Phase 8a slot-safe binding, issue #138)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/138
DefaultDependencies=no
After=var.mount
Before=local-fs.target umount.target
Conflicts=umount.target

[Mount]
What=/var/root
Where=/root
Type=none
Options=bind
TimeoutSec=10s

[Install]
WantedBy=local-fs.target
