# appliance/mkosi.extra/etc/systemd/system/spatiumddi-supervisor-watchdog.service
[Unit]
Description=SpatiumDDI external supervisor watchdog (#170 Wave E)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/170
# Don't run before docker.service is up — we need ``docker restart``
# to work; require + after dependency on the daemon. The container
# itself can be down (that's the case the watchdog cares about); we
# only need the daemon socket.
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-supervisor-watchdog
# The script itself never blocks for long (1 stat + maybe 1
# ``docker restart``), but cap with a sensible upper bound so a
# wedged docker daemon doesn't pile up overlapping invocations from
# the timer.
TimeoutStartSec=60
# Don't fail the unit on script exit != 0 — the script returns the
# underlying error code on docker-restart failures so the operator
# can see them in journalctl, but a transient docker hiccup is not
# a reason to drop the timer schedule.
SuccessExitStatus=1
