# appliance/mkosi.extra/etc/systemd/system/spatiumddi-verbose-boot-reload.service
[Unit]
Description=Apply pending SpatiumDDI verbose-boot toggle (grubenv spatium_verbose)
# The runner reads the trigger under /var/lib/spatiumddi/release-state and runs
# grub-editenv against the grubenv on the ESP (/boot/efi). RequiresMountsFor
# both PULLS IN and orders-after the mount units covering those paths — robust
# to the ESP mount-unit's generated name (vs a bare After=boot-efi.mount, which
# only orders if that unit happens to run). Without /boot/efi mounted the runner
# can't find grubenv.
RequiresMountsFor=/boot/efi /var/lib/spatiumddi/release-state
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-verbose-boot-reload
# No Install section — only invoked by the matching .path unit, never
# directly enabled at boot.
