# appliance/mkosi.extra/etc/systemd/system/spatiumddi-slot-rollback.service
[Unit]
Description=Apply pending SpatiumDDI slot rollback (Phase 8c-3)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/138
# No docker dependency — rollback just edits grubenv, no container
# state involved.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-slot-rollback
# No Install section — only invoked by the matching .path unit, never
# directly enabled at boot.
