# appliance/mkosi.extra/etc/systemd/system/spatiumddi-slot-set-next-boot.service
[Unit]
Description=Apply pending SpatiumDDI slot-set-next-boot trigger (Fleet UI)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/138

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-slot-set-next-boot
# No Install section — only invoked by the matching .path unit.
