# appliance/mkosi.extra/etc/systemd/system/spatiumddi-snmp-reload.service
[Unit]
Description=Apply pending SpatiumDDI snmpd config (issue #153)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/153

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-snmp-reload
# No Install section — only invoked by the matching .path unit, never
# directly enabled at boot.
