# appliance/mkosi.extra/etc/systemd/system/spatiumddi-reboot-agent.service
[Unit]
Description=Reboot the appliance on operator-triggered fleet request (Phase 8f-8)
Documentation=https://github.com/spatiumddi/spatiumddi/issues/138
# The reboot intent is operator-driven via the control plane's Fleet
# view — once the path unit fires we just need to actually reboot,
# no other host state is involved.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-reboot-agent
# No Install section — only invoked by the matching .path unit, never
# directly enabled at boot.
