# appliance/mkosi.extra/etc/systemd/system/spatiumddi-update.service
[Unit]
Description=Apply pending SpatiumDDI release (Phase 4c)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spatiumddi-update
# No Install section — this service is only invoked by the matching
# .path unit, never directly enabled at boot.
