# appliance/mkosi.extra/etc/systemd/system/spatium-install.service
[Unit]
Description=SpatiumDDI installer wizard (runs when boot cmdline has spatium-mode=install)
ConditionKernelCommandLine=spatium-mode=install
# NOTE: do NOT use `Conflicts=getty@tty1.service` here. systemd
# evaluates Conflicts= at unit-load time even when the unit's
# ConditionKernelCommandLine prevents activation — so on the
# *installed* system (no spatium-mode= flag), getty@tty1 ends up
# in ConflictedBy with this unit and never auto-starts. VGA
# console then has no login prompt. The stop is done as
# ExecStartPre below instead — only fires when this unit
# actually activates on the live ISO.
After=getty@tty1.service systemd-user-sessions.service plymouth-quit.service
After=systemd-networkd.service

[Service]
Type=idle
ExecStartPre=-/usr/bin/systemctl stop getty@tty1.service
ExecStart=/usr/local/bin/spatium-install
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
KillMode=process
IgnoreSIGPIPE=no
SendSIGHUP=yes
# DO NOT auto-restart. The wizard handles its own failure path
# (drops to shell with an error). Restarting would loop the
# operator back to question 1, which is what they hit before
# the kernel-console-spam fix landed.
Restart=no

[Install]
WantedBy=multi-user.target
