# appliance/mkosi.extra/etc/systemd/system/spatiumddi-supervisor-watchdog.timer
[Unit]
Description=Run SpatiumDDI external supervisor watchdog every 2 min
Documentation=https://github.com/spatiumddi/spatiumddi/issues/170

[Timer]
# Fire 60 s after boot — gives the supervisor's start-period
# (Dockerfile HEALTHCHECK has --start-period=15s, plus first heartbeat
# round-trip) time to land before the first liveness check.
OnBootSec=60s
# Every 2 min after that. Cadence is half the supervisor's heartbeat
# interval (60 s) so we catch a wedge within one missed loop +
# margin; bumping below 2 min adds churn without recovery-time gain
# since the script's STALE_THRESHOLD_S is 5 min anyway.
OnUnitActiveSec=2min
# Don't accumulate missed runs across a host pause/resume — we only
# care about the current state of the supervisor.
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target
