# appliance/mkosi.extra/etc/systemd/system/nftables.service.d/wait-for-etc-overlay.conf
# nftables.service ships with DefaultDependencies=no + Before=network-pre.target
# so it runs very early — before etc.mount has overlaid /var/persist/etc onto
# /etc. The result: the in-kernel ruleset is loaded from the IMAGE BASELINE
# /etc/nftables.conf, missing any operator edits that landed in the upper
# (/var/persist/etc/nftables.conf).
#
# Wait for etc.mount + spatium-etc-reconcile so the file `nft -f` reads is
# the post-overlay merged view. Before=network-pre.target stays in the unit
# itself, so any network-dependent service (systemd-networkd, k3s, …)
# still waits for the firewall to be live before it starts handling
# traffic.
[Unit]
After=etc.mount spatium-etc-reconcile.service
Wants=etc.mount spatium-etc-reconcile.service
