# appliance/mkosi.extra/etc/systemd/system/NetworkManager-wait-online.service.d/timeout.conf
# Drop NM-wait-online's wait from "wait for every profile to come
# fully online, default 30 s" to "wait up to 10 s for any
# connectivity, then continue regardless". A live appliance might
# legitimately boot on an isolated network for evaluation — blocking
# boot for 30 s on every isolated-network boot is bad UX. Real
# services (the docker stack, the firstboot orchestrator) handle
# network coming up late by themselves.
#
# Mirrors the equivalent drop-in we used to ship for
# systemd-networkd-wait-online before the NetworkManager switch.
[Service]
ExecStart=
ExecStart=/usr/bin/nm-online --quiet --timeout=10
