# appliance/mkosi.extra/etc/systemd/system/systemd-networkd-wait-online.service.d/timeout.conf
# Drop the wait-online block from "wait for EVERY configured
# interface, indefinitely" (the default) to "wait up to 10 s for
# ANY interface, then continue regardless". A live appliance might
# legitimately boot on an isolated network for evaluation — blocking
# boot for 2 minutes on every isolated-network boot is bad UX. Real
# services (the docker stack, the firstboot orchestrator) handle
# network coming up late by themselves.
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --timeout=10 --any
