# appliance/mkosi.extra/etc/systemd/system/serial-getty@ttyS0.service.d/no-loop.conf
# When the appliance has `console=ttyS0` in the kernel cmdline but
# no actual serial-port hardware (typical for Proxmox / KVM /
# VirtualBox VMs without an explicit serial device), agetty starts,
# can't open ttyS0 meaningfully, and exits — systemd then restarts
# it (default Restart=always) and we get an audit-spamming + console-
# slow restart loop every ~10 seconds forever.
#
# Cap to 3 failures within 30 seconds, then stop trying. Operators
# with a real serial port can `systemctl reset-failed
# serial-getty@ttyS0 && systemctl start serial-getty@ttyS0` to
# re-enable it on demand.
[Unit]
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
RestartSec=10s
