# appliance/mkosi.extra/etc/systemd/system/serial-getty@ttyS0.service.d/spatium-console-gate.conf
# #393 — gate serial-getty@ttyS0 on the appliance console mode, the serial
# mirror of getty@tty1.service.d/spatium-console-gate.conf.
#
# systemd-getty-generator auto-wants serial-getty@ttyS0 from the kernel
# cmdline `console=ttyS0,115200n8`. The serial console dashboard
# (spatium-console@ttyS0.service) runs when spatium-console=off is NOT set
# and owns ttyS0. This drop-in makes serial-getty its mirror image — it runs
# ONLY when spatium-console=off — so the two are mutually exclusive by
# condition and never race over ttyS0:
#   * dashboard / verbose_dashboard mode → this condition fails → the serial
#     dashboard owns ttyS0.
#   * text_console mode → the serial dashboard's condition fails, this passes
#     → a normal serial login.
# Without it (and without the now-removed Conflicts=), serial-getty would
# race the dashboard for ttyS0 in dashboard mode. NB: only validated on a box
# with a real serial port — the dev VM has none.
[Unit]
ConditionKernelCommandLine=spatium-console=off
